# scripts/systemd/pve-expire-lab.service
[Unit]
Description=PVE Lab VM Expiration
Documentation=file:///opt/pve-home/docs/VM-LIFECYCLE.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/pve-home/scripts/lifecycle/expire-lab-vms.sh --force
User=root
WorkingDirectory=/opt/pve-home
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pve-expire-lab

# Securite
ProtectSystem=strict
ReadWritePaths=/var/log
PrivateTmp=true
NoNewPrivileges=true
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectKernelLogs=true
RestrictRealtime=true
RestrictNamespaces=true
