# scripts/systemd/pve-cleanup-snapshots.service
[Unit]
Description=PVE Snapshot Cleanup
Documentation=file:///opt/pve-home/docs/VM-LIFECYCLE.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/pve-home/scripts/lifecycle/cleanup-snapshots.sh --force
User=root
WorkingDirectory=/opt/pve-home
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pve-cleanup-snapshots

# Securite
ProtectSystem=strict
ReadWritePaths=/var/log
PrivateTmp=true
NoNewPrivileges=true
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectKernelLogs=true
RestrictRealtime=true
RestrictNamespaces=true
